# crontab
* * * * * cd $PROJECT_ROOT_PATH && coursedog_env/bin/python manage.py runcrons
