# crontab
* * * * * cd $PROJECT_ROOT_PATH && python manage.py runcrons
